# systemd/twitse_client.service.template
[Unit]
Description=twitse client
Wants=network-online.target
After=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=<git root>
ExecStart=<git root>/bin/twitse_client --id <name>
Restart=always

[Install]
WantedBy=multi-user.target
